# tranquil-inbox-ward.service.template
[Unit]
Description=Tranquil Inbox Ward - Enhanced Email Classification Service
After=network.target ollama.service
Wants=network.target

[Service]
Type=simple
WorkingDirectory={INSTALL_DIR}
Environment="PATH={INSTALL_DIR}/{VENV_NAME}/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile={INSTALL_DIR}/.env

# 使用虚拟环境中的Python
ExecStart={INSTALL_DIR}/{VENV_NAME}/bin/gunicorn -w 4 -b 0.0.0.0:8501 --access-logfile /var/log/tranquil-inbox-ward/access.log --error-logfile /var/log/tranquil-inbox-ward/error.log app:app

# 重启策略
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=5

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log/tranquil-inbox-ward/
ProtectHome=true

# 日志设置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tranquil-inbox-ward

[Install]
WantedBy=multi-user.target
